# roles/defaults/files/update.service
[Unit]
Description=Update using pacman

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Syu --noconfirm
